# configs/com.claude.auto-attach.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.auto-attach</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/teample.casper/.claude/scripts/auto-attach.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/teample.casper</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/teample.casper/.claude/logs/auto-attach.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/teample.casper/.claude/logs/auto-attach.stderr.log</string>
</dict>
</plist>
